# files/tabouret.service
[Unit]
Description=Tabouret service

[Service]
ExecStart=/usr/bin/java -Dport=3003 -jar /home/ubuntu/tabouret-0.1.0-standalone.jar
